# templates/ucarp.debian.conf.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          ${:_buildout_section_name_}
# Required-Start:    $network $syslog $named
# Required-Stop:     $network $syslog $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start/stop ucarp daemon
# Description: This is a buildout generated file, do not edit directly.
#              It should be copied/symlinked in /etc/init.d and activated
#              Edit instead: ${:input}
#              using the update-rc.d script:
#              cp ${initscript:output} /etc/init.d/${initscript:name}
#              update-rc.d ${initscript:name} defaults
### END INIT INFO

# The Application Address that will failover
VIRTUAL_ADDRESS=${:vip_address}
#VIRTUAL_BROADCAST=<broadcast_ip>
#VIRTUAL_NETMASK=255.255.255.0

# The Maintanence Address of the local machine
SOURCE_ADDRESS=${:source_address}

# Instance id
# any number from 1 to 255
# Master and Backup need to be the same
INSTANCE_ID=${:uid}

# Location of the ucarp executable
UCARP_EXEC=/usr/sbin/ucarp
IFCONFIG=/sbin/ifconfig

# Location of the ucarp pid file
UCARP_PIDFILE=/var/run/ucarp.pid

# ucarp options
# -z run shudown script on exit
# -P force perferred master
# -B run in backgroud (DO NOT USE, it will change the pid (+1))
OPTIONS="{:options}"

# ucarp base, lower number will be perfered master
# set to same to have master stay alive as long as possible
UCARP_BASE=1

# ucarp skew, lower number will be perfered master (it is additive with UCARP_BASE)
UCARP_SKEW=${:skew}

# ucarp deadratio, number of sampling times to consider the daemon as dead
UCARP_DEADRATIO=4

# Interface for Ipaddress
INTERFACE=eth0

# Password so servers can trust who they are talking to
PASSWORD=${:password}

# The start and stop scripts
START_SCRIPT=${vip-up:output}
STOP_SCRIPT=${vip-down:output}

# Daemon Should start with this options.
# ucarp -z -B -b 2 -k 200 -r 4 --interface=eth0 --srcip=193.205.246.18 --vhid=1 --pass=zxcvbn --addr=193.205.246.16 --upscript=/etc/network/if-up.d/vip-up.sh --downscript=/etc/network/if-down.d/vip-down.sh

# DO NOT EDIT BELOW THIS LINE!
UCARP_OPTS="$OPTIONS --advbase=$UCARP_BASE --advskew=$UCARP_SKEW --deadratio=$UCARP_DEADRATIO --interface=$INTERFACE --srcip=$SOURCE_ADDRESS --vhid=$INSTANCE_ID --pass=$PASSWORD --addr=$VIRTUAL_ADDRESS --upscript=$START_SCRIPT --downscript=$STOP_SCRIPT"

case "$1" in
start)
  echo "Starting UCARP"
  start-stop-daemon --start --background --make-pidfile --pidfile ${UCARP_PIDFILE} --quiet --startas ${UCARP_EXEC} -- ${UCARP_OPTS}
  echo "... done"
  ;;
stop)
  echo "Stopping UCARP"
  start-stop-daemon --stop --pidfile ${UCARP_PIDFILE}
  echo "... done"
  ;;
*)
  echo "Usage: $0 {start|stop}" >&2
  exit 1
  ;;
esac

exit 0
